# Setup/Supervisors/moontracker_ui.conf
[program:moontracker_ui]
command=python3 main.py
directory=/home/pi/Automated_Earth_Moon_Earth_Communication/MoonTracker
autostart=true
autorestart=true
priority=3
user=pi
