# Setup/Supervisors/pigpiod.conf
[program:pigpiod]
command=pigpiod -l -g
autostart=true
autorestart=true
priority=1
user=root
